# skills/scheduled-task-validator/com.sam.validator.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sam.validator.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/samaguiar/Documents/Projects/Skills/scheduled-task-validator/validate.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>11</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key>
    <string>/Users/samaguiar/Documents/Codex/validator-runs/_launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/samaguiar/Documents/Codex/validator-runs/_launchd-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/samaguiar/Documents/Projects/Skills/scheduled-task-validator</string>
    <key>KeepAlive</key><false/>
</dict>
</plist>